# files/etc/cron.d/cloak
@daily root /bin/update-pki.sh
* * * * * root /bin/refresh-crls.sh
* * 1,20 * * /usr/bin/certbot renew --force-renew --deploy-hook /bin/update-certificate.sh -q > /dev/null 2>&1
