# files/etc/cron.d/cloak
@daily root /bin/update-pki.sh
* * * * * root /bin/refresh-crls.sh
