# files/etc/cron.d/cloak
@daily cloak /home/cloak/bin/update-pki.sh
* * * * * cloak /home/cloak/bin/refresh-crls.sh
