# files/etc/cron.d/cloak
#min	hour	day/month	month	day/week
@daily root /usr/bin/update-pki.sh
*/5 * * * * root /usr/bin/refresh-crls.sh
10 0 1,8,15,22 * * root /usr/bin/renew-cert.sh
