# to_copy/etc/cron.d/cloak
@daily root /bin/update-pki.sh
*/5 * * * * root /bin/refresh-crls.sh
10 0 1,8,15,22 * * root /usr/bin/renew-cert.sh
*/10 * * * * root /usr/bin/boot-expired-revoked-certs.sh

