# to_copy/etc/cron.d/cloak
@daily root /bin/update-pki.sh
*/5 * * * * root /bin/refresh-crls.sh
0 0 1,15 * * root /usr/bin/renew-cert.sh
